# kde/xdg-desktop-portal-kde/src/xdg-desktop-portal-kde-6.2.2/data/plasma-xdg-desktop-portal-kde.service.in
[Unit]
Description=Xdg Desktop Portal For KDE
PartOf=graphical-session.target
After=plasma-core.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/xdg-desktop-portal-kde
BusName=org.freedesktop.impl.portal.desktop.kde
Slice=session.slice
# Auto restart is off as we are DBus activated
Restart=no
