# tde/tde-core/tde-tdebase/tdm.service
[Unit]
Description=TDE Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=/opt/trinity/bin/tdm

[Install]
Alias=display-manager.service
