# gnupg/dirmngr@.socket
[Unit]
Description=GnuPG network certificate management daemon for %f
Documentation=man:dirmngr(8)

[Socket]
ListenStream=%f/S.dirmngr
SocketMode=0600
