# gnupg/dirmngr@.service
[Unit]
Description=GnuPG network certificate management daemon for %f
Documentation=man:dirmngr(8)
Requires=dirmngr@%i.socket

[Service]
ExecStart=/usr/bin/dirmngr --homedir %f --supervised
ExecReload=/usr/bin/gpgconf --homedir %f --reload dirmngr
