# gnupg/gpg-agent-extra@.socket
[Unit]
Description=GnuPG cryptographic agent and passphrase cache (restricted) for %f
Documentation=man:gpg-agent(1)

[Socket]
ListenStream=%f/S.gpg-agent.extra
FileDescriptorName=extra
Service=gpg-agent@%i.service
SocketMode=0600
