# gnupg/gpg-agent-ssh@.socket
[Unit]
Description=GnuPG cryptographic agent (ssh-agent emulation) for %f
Documentation=man:gpg-agent(1) man:ssh-add(1) man:ssh-agent(1) man:ssh(1)
ConditionPathExists=%f
ConditionPathIsReadWrite=%f

[Socket]
ListenStream=%f/S.gpg-agent.ssh
FileDescriptorName=ssh
Service=gpg-agent@%i.service
SocketMode=0600
