# gnupg/gpg-agent@.service
[Unit]
Description=GnuPG cryptographic agent and passphrase cache for %f
Documentation=man:gpg-agent(1)
Requires=gpg-agent@%i.socket

[Service]
ExecStart=/usr/bin/gpg-agent --homedir %f --supervised
ExecReload=/usr/bin/gpgconf --homedir %f --reload gpg-agent
