# gnupg/keyboxd@.socket
[Unit]
Description=GnuPG public key management service for %f
ConditionPathExists=%f
ConditionPathIsReadWrite=%f

[Socket]
ListenStream=%f/S.keyboxd
FileDescriptorName=std
SocketMode=0600
