# gnupg/keyboxd.socket
[Unit]
Description=GnuPG public key management service

[Socket]
ListenStream=%t/gnupg/S.keyboxd
FileDescriptorName=std
SocketMode=0600
DirectoryMode=0700
