# gnupg/keyboxd@.service
[Unit]
Description=GnuPG public key management service for %f
Requires=keyboxd@%i.socket

[Service]
ExecStart=/usr/lib/gnupg/keyboxd --homedir %f --supervised
ExecReload=/usr/bin/gpgconf --homedir %f --reload keyboxd
